# server-configs/chromium-rdp-watchdog/chromium-rdp-watchdog.timer
[Unit]
Description=Check RDP Chromium CDP :9223 every 30 seconds

[Timer]
OnBootSec=1min
OnUnitActiveSec=30s
AccuracySec=5s
Persistent=true

[Install]
WantedBy=timers.target
